# files/master/etc/uci-defaults/91-system-crontab
cat << "EOF" >> /etc/crontabs/root
*/2 * * * * /usr/bin/wireguard_watchdog
EOF

uci set system.@system[0].cronloglevel="9"
uci commit system